# scripts/couchdb-replicate.service
[Unit]
Description=Replicate CouchDB obsidianvault <-> obsidianknowledge

[Service]
Type=oneshot
ExecStart=/opt/HomeLab/scripts/couchdb-replicate.sh
