# scripts/couchdb-replicate.timer
[Unit]
Description=Run CouchDB replication every minute

[Timer]
OnCalendar=*-*-* *:*:30
Persistent=true

[Install]
WantedBy=timers.target
