# scripts/goodnotes-ingest.timer
[Unit]
Description=Run GoodNotes ingest every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target